# resource/refresh.cron
0 0 * * * /bin/systemctl restart incron.service
